# nginx_api.service
# copy this file to /lib/systemd/system/

[Unit]
Description=NGINX API
After=network-online.target

[Service]
Restart=on-failure
WorkingDirectory=/home/harmonic/api/
ExecStart=/usr/bin/node /home/harmonic/api/nginx_api.js

[Install]
WantedBy=multi-user.target
